# root/etc/update-motd.d/10-retro-home-help-text
#!/bin/sh

printf "\n"
printf " 🐙 Code:        https://github.com/Sunderland93/ReplayOS\n"
printf " 🐛 Bugs:        https://github.com/Sunderland93/ReplayOS/issues\n"
printf "\n"
